# crontab
# START CRON JOB
# * * * * * python3 /app/pingapi.py >> /proc/1/fd/1 2>/proc/1/fd/2
* * * * * python3 /app/pingapi.py >> /var/log/cron.log 2>&1
# END CRON JOB